# plasma-login-manager/packages/plasma-login-manager/postinst
#!/bin/bash

if [ ! -L /etc/systemd/system/display-manager.service ]; then
systemctl -f enable plasmalogin
systemctl set-default graphical.target
fi
