# plasma-login-manager/packages/plasma-login-manager/prerm
#!/bin/bash

systemctl -s stop plasmalogin
systemctl -f disable plasmalogin
systemctl set-default multi-user.target
